# ops/rnbp-api.service
[Unit]
Description=RNBP API Server
After=network.target

[Service]
Type=simple
User=prod
Group=prod
WorkingDirectory=/opt/rnbp/repo/apps/api
ExecStart=/home/prod/.nvm/versions/node/v20.19.4/bin/node dist/index.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rnbp-api

# Environment
EnvironmentFile=/opt/rnbp/.env

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/opt/rnbp/uploads
PrivateTmp=true

[Install]
WantedBy=multi-user.target
